# ops/systemd/aurora-monitor.timer
[Unit]
Description=Run Aurora production monitor every five minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
RandomizedDelaySec=30s
Persistent=true
Unit=aurora-monitor.service

[Install]
WantedBy=timers.target
