# ops/systemd/aurora-monitor.service
[Unit]
Description=Verify Aurora production container health and resource headroom
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/aurora/monitor.env
ExecStart=/opt/aurora/ops/monitor-container.sh
User=root
Group=root
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
